# rc.local
#!/bin/bash
systemctl restart libvirtd.service
virsh start win10-z-clone
